# deployment/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=/usr/local/bin/gunicorn vacayo.wsgi -k gevent -w 2 -t 120 -b 127.0.0.1:8000 --chdir=/opt/app

[program:nginx]
command=/usr/sbin/nginx
